# quiz.service
[Unit]
Description=Gunicorn instance to serve quiz
After=network.target

[Service]
User=raphaeluziel
Group=www-data
WorkingDirectory=/home/raphaeluziel/quiz
Environment="PATH=/home/raphaeluziel/quiz/quizenv/bin"
ExecStart=/home/raphaeluziel/quiz/quizenv/bin/gunicorn --workers 1 --bind unix:quiz.sock -m 007 wsgi:app


[Install]
WantedBy=multi-user.target
